# run-image/service/s6-svscan-log/run
#!/bin/execlineb -P

if { install -d -m 0750 -o log -g log /run/uncaught-logs }
redirfd -w 2 /dev/console
redirfd -w 1 /dev/null
redirfd -rnb 0 fifo
umask 027
s6-setuidgid log
s6-log -bpd3 -- !zstd t /run/uncaught-logs
