# run-image/service/s6-linux-init-runleveld/run
#!/bin/execlineb -P

fdmove -c 2 1
fdmove 1 3
s6-ipcserver -1 -a 0700 -c 1 -- s
s6-sudod -0 -1 -2 -t 30000 --
"/etc/rc"/scripts/runlevel
