# run-image/service/s6-linux-init-shutdownd/run
#!/bin/execlineb -P

s6-linux-init-shutdownd -c "/etc/rc" -g 2000
